# debian/rpi-imager-embedded.rpi-imager-eglfs.service
[Unit]
Description=Raspberry Pi Imager (Embedded EGLFS Mode)
Documentation=https://github.com/raspberrypi/rpi-imager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rpi-imager-embedded
Restart=on-failure
RestartSec=5
User=root
Group=root
Environment=HOME=/root
Environment=XDG_RUNTIME_DIR=/run/user/0

# EGLFS specific environment
Environment=QT_QPA_PLATFORM=eglfs
Environment=QT_QPA_EGLFS_INTEGRATION=eglfs_kms
Environment=QT_QPA_EGLFS_KMS_DEVICE=/dev/dri/card0
Environment=QT_QPA_EGLFS_ALWAYS_SET_MODE=1
Environment=QT_QUICK_CONTROLS_STYLE=Material
Environment=QT_OPENGL=es2

[Install]
WantedBy=multi-user.target 